# templates/mysqlchk.socket.epp
[Unit]
Description=Mysqlchk socket

[Socket]
ListenStream=<%= $galera::status_port %>
Accept=yes

[Install]
WantedBy=sockets.target
