# templates/mysqlchk.service.epp
[Unit]
Description=Mysqlchk service
After=network.target network-online.target
Wants=network-online.target

[Service]
User=clustercheck
Group=clustercheck
StandardInput=socket
ExecStart=/usr/local/bin/clustercheck
